# custom/fs-skeleton/etc/crontabs/root
59 * * * * /usr/local/bin/hourly-paused-jpeg2mp4-converter.sh

